# conf/systemd.path
[Unit]
Description="Lauch 'zola build' on content change"

[Path]
PathModified=__INSTAL_DIR__/content
PathModified=__INSTAL_DIR__/static
PathModified=__INSTAL_DIR__/templates
PathModified=__INSTAL_DIR__/config.toml
PathModified=__INSTAL_DIR__/themes
Unit=__APP__.service

[Install]
WantedBy=multi-user.target